# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /roo/root/.pyenv/versions/Deploy-PlotTP/bin/uwsgi --ini /srv/Deploy-PlotTP/.config/uwsgi/uwsgi-app.ini